# .config/systemd/user/dufs.service
[Unit]
Description=Dufs File Server
After=network.target

[Service]
ExecStart=/usr/bin/dufs -A
Restart=always

[Install]
WantedBy=default.target
