# .config/systemd/user/wallpaper-minute.service
[Unit]
Description=Set wallpaper by minute (swww)

[Service]
Type=oneshot
ExecStart=%h/.config/niri/wallpaper-minute.sh
